# templates/rclone-backup.service.erb
# Puppet managed file
[Unit]
OnFailure=alert-slack@%N.service
OnFailure=+/bin/bash -c "/usr/bin/echo '<%= @name %>-backup started' | /usr/bin/mailx -s '<%= @name %>-backup failed' <%= @email %>"

[Service]
Type=oneshot
User=<%= @user %>
Group=<%= @group %>
ExecStartPre=+/usr/local/bin/slack_msg --channel notifications --icon robot --message '<%= @name %>-backup started'
ExecStartPre=+/bin/bash -c "/usr/bin/echo '<%= @name %>-backup started' | /usr/bin/mailx -s '<%= @name %>-backup started' <%= @email %>"
<% if @pre_rclone != nil -%>
<% @pre_rclone.each do |pre_exec| -%>
ExecStartPre=<%= pre_exec %>
<% end -%>
<% end -%>
EnvironmentFile=-/tmp/systemd_<%= @user %>_rclone.env
ExecStart=/usr/bin/rclone <%= @command %> <%= @rclone_opts %> <%= @src %> <%= @dst %> --log-file /var/log/rclone-backups/<%= @name %>-backup.log --stats=0 --stats-one-line -v
ExecStartPost=+/bin/bash -c "/usr/bin/echo '<%= @name %>-backup completed' | /usr/bin/mailx -s '<%= @name %>-backup completed' -a '/var/log/rclone-backups/<%= @name %>-backup.log' <%= @email %>"
<% if @post_rclone != nil -%>
<% @post_rclone.each do |post_exec| -%>
ExecStartPost=<%= post_exec %>
<% end -%>
<% end -%>
ExecStartPost=+/usr/local/bin/slack_msg --channel notifications --icon robot --message '<%= @name %>-backup completed'
